# claude-voice-assistant.service
[Unit]
Description=Toris Agent — Claude voice bot for Telegram
After=network.target

[Service]
Type=simple
User=claude
Group=claude
WorkingDirectory=/opt/toris-agent
EnvironmentFile=/etc/toris-agent/toris-agent.env
ExecStart=/opt/toris-agent/venv/bin/python /opt/toris-agent/bot.py
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/toris-agent/sandbox /var/lib/toris-agent/state

[Install]
WantedBy=multi-user.target
